# ser2net/rootfs/etc/services.d/ser2net/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Ser2net
# Runs Ser2net
# ==============================================================================

exec ser2net -l -c /etc/ser2net/ser2net.yaml